# deploy/bot_rabbitmq.service
[Unit]
Description=Telegram bot rabbitmq consumer
After=multi-user.target

[Service]
Type=idle
User=user
Group=user
WorkingDirectory=/home/user/market/backend/tg_bot/src
ExecStart=poetry run python3.11 start_consuming.py
Restart=always

[Install]
WantedBy=multi-user.target
